# system/dashboard.service
[Unit]
Description=PTP Dashboard Web Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/pi5-ptp-grandmaster/backend
ExecStart=/usr/bin/python3 app.py

# Restart behavior
Restart=always
RestartSec=3
KillMode=process
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Environment (centralized config)
EnvironmentFile=/opt/pi5-ptp-grandmaster/system/environment.conf

# Run as root (or create a dedicated user)
User=root

[Install]
WantedBy=multi-user.target
