# system/logger.service
[Unit]
Description=PTP Data Logging Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/pi5-ptp-grandmaster/backend/logging_service.py
WorkingDirectory=/opt/pi5-ptp-grandmaster/backend
Restart=always
User=root

[Install]
WantedBy=multi-user.target
